# examples/systemd/nextcloud/extension-XX_nextcloud_cron.service
[Unit]
Description=Nextcloud Cronjob

[Service]
ExecStart=podman exec --user www-data XX_nextcloud php -f /var/www/html/cron.php
KillMode=process
